# deploy/automation.service
[Unit]
Description=MY_PROJECT Automation (Threads + YouTube + Article)
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
# Sesuaikan WorkingDirectory, path bun, dan User saat deploy.
WorkingDirectory=/opt/automation
ExecStart=/root/.bun/bin/bun run src/all.ts
EnvironmentFile=/opt/automation/.env
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
